# home/dot_config/systemd/user/deck-backup.service
[Unit]
Description="Deck backup"

[Service]
ExecStart=%h/.local/bin/deck-backup -f
