# home/dot_config/systemd/user/deck-backup.service
[Unit]
Description="Deck backup"

[Service]
ExecStart=/home/deck/.local/bin/deck-backup -f
