# home/dot_config/systemd/user/nix-clean.service.tmpl
{{- if and (eq .chezmoi.os "linux") (not .ephemeral) .use_nix (lookPath "nix") -}}
[Unit]
Description=Nix Cleanup - Garbage collection and store optimization
Documentation=man:nix-collect-garbage(1)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart={{ .chezmoi.destDir }}/bin/nix-clean
StandardOutput=journal
StandardError=journal

# Resource limits (low priority)
Nice=15

# I/O priority (best-effort class, priority 7 = lowest)
IOSchedulingClass=best-effort
IOSchedulingPriority=7

# Track I/O usage
IOAccounting=true

# Sandboxing for security
PrivateTmp=true
NoNewPrivileges=true

# Environment
Environment=PATH={{ .chezmoi.destDir }}/bin:{{ .chezmoi.destDir }}/.nix-profile/bin:{{ .chezmoi.destDir }}/.local/state/nix/profiles/profile/bin:/nix/var/nix/profiles/default/bin:{{ .chezmoi.destDir }}/.brew/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin

[Install]
WantedBy=default.target
{{- end -}}
