# home/dot_config/systemd/user/retrodeck-upgrade.timer.tmpl
{{ if .retrodeck }}
[Unit]
Description=Weekly RetroDECK upgrade

[Timer]
OnCalendar=weekly
Persistent=true
Unit=retrodeck-upgrade.service

[Install]
WantedBy=timers.target
{{ end }}
